# files/etc/rc.local
#!/bin/ash

/bin/charge-set

/TR750/battery/LowPowerCheck.sh &

exit 0
